# systemd/lora-service.service
[Unit]
Description=SOS-GUIDE LoRa Mesh Service
Documentation=https://github.com/sos-guide
After=network.target
Wants=network.target

# Ne démarre que si le flag enableLoRa=true dans config.json
# (vérifié dans ExecStartPre)
ConditionPathExists=/var/www/sos-guide/data/config.json

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/local/bin

# Vérifier que LoRa est activé avant de démarrer
ExecStartPre=/bin/bash -c \
    "jq -e '.enableLoRa == true' /var/www/sos-guide/data/config.json > /dev/null \
     || { echo 'LoRa désactivé dans config.json — service non démarré'; exit 1; }"

ExecStart=/usr/bin/python3 /usr/local/bin/lora-service.py

# Redémarrage automatique (sauf arrêt volontaire)
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Ressources
MemoryMax=64M
CPUQuota=20%
Nice=10

# Sécurité systemd
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log /etc/sos-guide /var/lib/sos-guide
PrivateTmp=yes

# Accès SPI requis pour SX1276
DeviceAllow=/dev/spidev0.0 rw
DeviceAllow=/dev/gpiomem rw
SupplementaryGroups=gpio spi

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lora-service

[Install]
WantedBy=multi-user.target
