# systemd/sos-guide-firstboot.service
[Unit]
Description=SOS-GUIDE First Boot Setup (STARTER mode)
After=network.target
Before=hostapd.service dnsmasq.service nginx.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/firstboot.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
